# infra/systemd/bbctl-docops-sync.timer
[Unit]
Description=Run bbctl docops sync nightly

[Timer]
OnCalendar=*-*-* 02:30:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
