# infra/systemd/bbctl-repo-sync.timer
[Unit]
Description=Run bbctl repo sync nightly

[Timer]
OnCalendar=*-*-* 02:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
